# Keychron
#!/sbin/openrc-run

# Keychron keyboard function key support for OpenRC systems (Tested with gentoo)

echo 2 > /sys/module/hid_apple/parameters/fnmode

exit

